# meta-mender-demo/recipes-mender/boot-script/files/boot-script.service
[Unit]
Description=Boot script to demo Mender OTA updates

[Service]
ExecStart=/usr/sbin/boot-script.sh

[Install]
WantedBy=multi-user.target
