# recipes-mender/mender/files/mender.service
[Unit]
Description=Mender OTA update service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/mender -daemon -log-level=debug -log-file=/var/log/mender.log
Restart=on-abort

[Install]
WantedBy=multi-user.target
