# meta-mender-core/recipes-mender/mender-client/files/mender-client-data-dir.service
[Unit]
Description=Mender persistent data dir
After=data.mount
ConditionPathExists=!/data/mender

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/bin/mkdir -p -m 0700 /data/mender

[Install]
WantedBy=mender-client.service
